# deploy/systemd/cloudflared-opencode-mcp-gateway.service
[Unit]
Description=Cloudflare Tunnel for OpenCode MCP Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gabriel
Group=gabriel
ExecStart=/tmp/cloudflared tunnel run --token eyJhIjoiNzFjM2VmMGI3ZGE5NTU5NDNjMTI5OTg4ZGI3MTAzNmIiLCJ0IjoiMjc1OTlhNGQtMWQ2Ni00ZGQwLWJjNjYtOGMyNGRkZDY1NWRkIiwicyI6Ik1Ub0l0K2hHQ00rNGdDbFZWa2JxM0NRY3lEMVVtL3RXTk5NY1VQRTBSWlE9In0=
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
