# deploy/systemd/opencode-mcp-gateway@.service
[Unit]
Description=OpenCode MCP Gateway instance %i
After=network-online.target opencode.service
Wants=network-online.target

[Service]
Type=simple
User=gabriel
Group=gabriel
WorkingDirectory=/home/gabriel/AI Projects:/opencode-mcp-gateway
EnvironmentFile=/home/gabriel/AI Projects:/opencode-mcp-gateway/.env.%i
ExecStart=/usr/bin/env python3 /home/gabriel/AI Projects:/opencode-mcp-gateway/main.py
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
